# package/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Fix permissions on helper scripts
        chmod 755 /usr/lib/dw-iface/dw-up.sh \
                  /usr/lib/dw-iface/dw-down.sh \
                  /usr/lib/dw-iface/dw-status.sh \
                  /usr/lib/dw-iface/dw-doctor.sh

        # Create config dir if not present
        mkdir -p /etc/dw-iface /run/dw-iface /var/log/dw-iface

        # Install example config if no config exists yet
        if [ ! -f /etc/dw-iface/dw-iface.conf ]; then
            cp /etc/dw-iface/dw-iface.conf.example /etc/dw-iface/dw-iface.conf
            echo "dw-iface: created /etc/dw-iface/dw-iface.conf — edit before running 'dw-iface up'"
        fi

        # Enable systemd service but do NOT start it (radio may not be connected)
        # Enable service only when systemd is the active init (not in Docker/chroot)
        if command -v systemctl >/dev/null 2>&1 \
                && systemctl is-system-running >/dev/null 2>&1; then
            systemctl daemon-reload 2>/dev/null || true
            systemctl enable dw-iface.service 2>/dev/null || true
            echo "dw-iface: service enabled — start with: sudo systemctl start dw-iface"
        fi
        ;;
esac

#DEBHELPER#
exit 0
